# crontab
#put this into "crontab -e" on server
1 * * * *      cp -r ~/Results /tmp/
